# data/ayatana-indicator-power.conf.in
description "Ayatana Indicator Power Service"

start on ayatana-indicator-services-start
stop on desktop-end or ayatana-indicator-services-end

respawn
respawn limit 2 10

exec @pkglibexecdir@/ayatana-indicator-power-service
