# data/ayatana-indicator-power.service.in
[Unit]
Description=Ayatana Indicator Power
PartOf=graphical-session.target
PartOf=ayatana-indicators.target

[Service]
ExecStart=@pkglibexecdir@/ayatana-indicator-power-service
Restart=on-failure

[Install]
WantedBy=ayatana-indicators.target
