# data/ayatana-indicator-power.service.in
[Unit]
Description=Ayatana Indicator Power
PartOf=graphical-session.target
After=ayatana-indicators-pre.target

[Service]
ExecStart=@pkglibexecdir@/ayatana-indicator-power-service
Restart=on-failure
